# supervisor.d/asterisk.conf
[program:asterisk]
command=/usr/sbin/asterisk -T -W -U asterisk -p -vvvdddf
autostart = true
startsecs = 80
startretries = 5
autorestart = true
priority=222
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stdout_events_enabled = true
stderr_events_enabled = true