# supervisor.d/mysql.conf
[program:mysql]
command=/usr/libexec/mysqld -u mysql
user=root
autostart=true
startsecs=30
autorestart=true
numprocs=1
